# examples/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude-opencode-bridge</string>
  <key>ProgramArguments</key>
  <array>
    <string>YOUR_BRIDGE_BIN</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>7411</string>
    <string>--session-store</string>
    <string>$HOME/.claude-opencode-bridge/sessions.json</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
</dict>
</plist>
